# s6-overlay/s6-rc.d/postgres/finish
#!/command/execlineb -S0

# Log shutdown
foreground { echo "PostgreSQL service finished" }

# Exit with the original exit code
exit ${1}
