# s6-overlay/s6-rc.d/init-services/up
/etc/s6-overlay/scripts/init-services.sh